# build/supervisor/web-server.conf
[program:web-server]
user=vagrant
directory = /vagrant/project/birdr/
command=python3 manage.py runserver 0:8000
autostart=true
stderr_logfile=/home/vagrant/logs/web-server.log
stdout_logfile=/home/vagrant/logs/web-server.log
